# cron/cron.daily/backup
#!/bin/sh

rsync -av \
      --files-from=/home/bbuccianti/.backup-dirs \
      --exclude-from=/home/bbuccianti/.exclude-dirs \
      /home/bbuccianti \
      backup@atom:backup/
