# cron/cron.daily/backup
#!/bin/sh

rsync -avr \
      --files-from=/home/bex/.backup-dirs \
      --exclude-from=/home/bex/.exclude-dirs \
      /home/bex \
      backup@atom:bex/
